# crontab
*/10 * * * * /home/zyin/rssfulltextgenerator/GenerateAll.py > ~/CronLog/rssfulltextgenerator.log 2>&1
